# bin/oowlish
#!/bin/bash

### BEGIN INIT INFO
# Provides:          Oowlish
# Short-Description: Oowlish
# Description:       Run Oowlish service
### END INIT INFO

case "$1" in
  start)
    echo "Starting Oowlish..."
    sudo -u ec2-user bash -c '/home/ec2-user/oowlish/bin/start-oowlish.sh'
    ;;
  stop)
    echo "Stopping Oowlish..."
    kill -9 `ps -ef |grep -i "python3 /home/ec2-user/oowlish/app.py" | grep -v grep | awk '{print $2}'`
    ;;
  *)
    echo "Usage: /etc/init.d/oowlish {start|stop}"
    exit 1
    ;;
esac

exit 0
